# .cron
28 6 10 * * sudo /usr/local/vesta/bin/v-backup-full-users
10 05 * * * sudo /usr/local/vesta/bin/v-backup-users && sudo /usr/local/vesta/bin/v-backup-daily-cloud-yandex
